# prune/docker-prune.timer.jinja
[Unit]
Description=Timer to automatically cleanup Docker images/ containers

[Timer]
OnCalendar={{ on_calendar }}
Unit=docker-prune.service

[Install]
WantedBy=basic.target
